# deploy/systemd/quantos-agent-daily-restart.service
[Unit]
Description=Daily restart of QuantOS Trading Agent (bounds unexplained memory growth)

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart quantos-agent.service
